# backend/deployment/arboris-async-processor.service
# Arboris Novel 异步分析处理器服务
#
# 使用方法：
# 1. 编辑本文件，替换 {{PROJECT_DIR}} 和 {{USER}} 为实际值
# 2. 复制到 /etc/systemd/system/: sudo cp arboris-async-processor.service /etc/systemd/system/
# 3. 重载 systemd: sudo systemctl daemon-reload
# 4. 启动服务: sudo systemctl start arboris-async-processor
# 5. 开机自启: sudo systemctl enable arboris-async-processor
# 6. 查看状态: sudo systemctl status arboris-async-processor
# 7. 查看日志: sudo journalctl -u arboris-async-processor -f

[Unit]
Description=Arboris Novel Async Analysis Processor
After=network.target arboris-api.service

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_DIR}}/backend
Environment="PATH={{PROJECT_DIR}}/backend/venv/bin"
ExecStart={{PROJECT_DIR}}/backend/venv/bin/python -m app.background_processor
Restart=always
RestartSec=10

# 日志
StandardOutput=append:/var/log/arboris/async-processor.log
StandardError=append:/var/log/arboris/async-processor-error.log

# 资源限制
LimitNOFILE=65536
MemoryLimit=2G

[Install]
WantedBy=multi-user.target
